# src/main/assembly/dist/bin/easy-bag-index.service
[Unit]
Description=EASY Bag Index Service

[Service]
ExecStart=/opt/dans.knaw.nl/easy-bag-index/bin/easy-bag-index-service-start

User=easy-bag-index
Group=easy-bag-index

[Install]
WantedBy=multi-user.target
