# src/main/assembly/dist/install/easy-bag-index.service
[Unit]
Description=EASY Bag Index Service

[Service]
ExecStart=/bin/java \
 -Dlogback.configurationFile=/etc/opt/dans.knaw.nl/easy-bag-index/logback-service.xml \
 -Dapp.home=/opt/dans.knaw.nl/easy-bag-index \
 -jar /opt/dans.knaw.nl/easy-bag-index/bin/easy-bag-index.jar run-service

User=easy-bag-index
Group=easy-bag-index

[Install]
WantedBy=multi-user.target
